# docker/bind9/rootfs/etc/s6-overlay/s6-rc.d/bind9/finish
#!/bin/bash
bash /scripts/docker-remove-bind.sh
echo "Running bind remove script"
