# docker/bind9/rootfs/etc/s6-overlay/s6-rc.d/bind9/run
#!/bin/bash
# run the bind install script
echo "Running bind install script"
eval /scripts/docker-install-bind.sh "${BCIP}" "${NSIP}"

